# deploy/addwii-bot-a.service
[Unit]
Description=addwii-pm Bot A（總監 AI 代理人）
After=network.target
Wants=network.target

[Service]
Type=simple
User=claude
WorkingDirectory=/opt/addwii-pm
Environment=PYTHONUTF8=1
EnvironmentFile=/opt/addwii-pm/.env
ExecStart=/opt/addwii-pm/.venv/bin/python3 main.py --role A
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=addwii-bot-a

[Install]
WantedBy=multi-user.target
